# launchd/com.sjr.appscope.snapshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Daily AppScope snapshot at 09:00.
  Install:
    cp launchd/com.sjr.appscope.snapshot.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.sjr.appscope.snapshot.plist
  Logs: /tmp/appscope-snapshot.log
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sjr.appscope.snapshot</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/npm</string>
    <string>run</string>
    <string>snapshot</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/sjr/Projects/appscope</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/appscope-snapshot.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/appscope-snapshot.log</string>
</dict>
</plist>
